# supervisord.conf
[supervisord]
nodaemon=true

[program:fakes3]
command=/usr/local/bin/fakes3 -r /fakes3_root -p 4569

[program:saboteur]
command=/usr/bin/saboteur-agent
